# portfolio-valentine-main/portfolio-valentine-prod.service
[Unit]
Description=Portfolio Valentine - PRODUCTION
Documentation=https://valentine-arnaly.com
After=network.target
Wants=portfolio-valentine-preprod.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/portfolio-valentine-prod/backend
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=10

# Variables d'environnement PRODUCTION
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=SITE_URL=https://valentine-arnaly.com
Environment=JWT_SECRET=production-super-secret-jwt-key-valentine-2025
Environment=SMTP_HOST=mail.infomaniak.com
Environment=SMTP_PORT=465
Environment=SMTP_USER=hello@buro.ooo
Environment=SMTP_PASS=Th1s-1s-Bur0
Environment=SMTP_FROM=hello@buro.ooo
Environment=DEV_SHOW_MAGIC_LINKS=false

# Logs avec préfixe PROD
StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfolio-valentine-prod

# Sécurité renforcée pour la production
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/www/portfolio-valentine-prod
ProtectHome=yes
ProtectKernelTunables=yes

[Install]
WantedBy=multi-user.target